# app-misc/heaviest/files/heaviest.initd
#!/sbin/openrc-run
# Copyright 1999-2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command=/usr/bin/heaviest-cpu-daemon
command_background="true"
pidfile="/run/heaviest.pid"
